# com.ambity.sentinel.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ambity.sentinel.watchdog</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/macbookpro/Downloads/bot_project/venv/bin/python3</string>
        <string>/Users/macbookpro/Downloads/bot_project/watchdog.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/macbookpro/Downloads/bot_project</string>

    <!-- Auto-restart watchdog if it crashes -->
    <key>KeepAlive</key>
    <true/>

    <!-- Start at login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Logs -->
    <key>StandardOutPath</key>
    <string>/Users/macbookpro/Downloads/bot_project/logs/watchdog_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/macbookpro/Downloads/bot_project/logs/watchdog_stderr.log</string>

    <!-- Throttle: if watchdog crashes repeatedly, wait 30s before restart -->
    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
